# zenbook-duo-session.service
[Unit]
Description=Zenbook Duo Session Daemon
# GNOME sets XDG_CURRENT_DESKTOP to values such as "GNOME" or "ubuntu:GNOME".
After=graphical-session.target gnome-session-initialized.target
PartOf=graphical-session.target

[Service]
Type=simple
# Only start under a GNOME graphical session (physical login uses graphical-session.target).
ExecCondition=/bin/sh -c 'case "${XDG_CURRENT_DESKTOP:-}" in *GNOME*) exit 0;; *) exit 1;; esac'
ExecStart=/opt/zenbook-duo-daemon/zenbook-duo-daemon session
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
